# roles/setup-buildslave/templates/buildbot.service
[Unit]
Description=Buildbot

[Service]
ExecStart=/buildslave/start_nd
Restart=always
User=buildslave

[Install]
WantedBy=multi-user.target
